# docker/apache2-supervisord.conf
[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && mkdir -p ${APACHE_LOCK_DIR} ${APACHE_LOG_DIR} ${APACHE_RUN_DIR} && /usr/sbin/apache2 -DFOREGROUND"

